# configs/stone-node.service
[Unit]
Description=Stone Node __STONE_NETWORK__ (Unified Setup + Full Node)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=__STONE_ROOT__
EnvironmentFile=-__STONE_ROOT__/.env
Environment=STONE_NETWORK=__STONE_NETWORK__
ExecStart=__STONE_PATH__/stone-setup
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stone-node-__STONE_NETWORK__
LimitNOFILE=65535

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
